# package/root/etc/update-motd.d/05-figlet
#!/bin/sh

test -x /usr/bin/figlet || exit 0

figlet $(hostname)
